# klipper/moonraker.ini
[program:moonraker]
user=john
command=/home/john/moonraker-env/bin/python /home/john/moonraker/moonraker/moonraker.py -l /var/log/klipper/moonraker.log -c /home/john/.config/moonraker.conf
environment=USER=john,HOME=/home/john
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes = 0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes = 0
autorestart=true
autostart=false
dependent_startup=true
dependent_startup_wait_for=klipper:running
